# openclaw-watchdog.service
[Unit]
Description=OpenClaw Agent Watchdog
After=network.target openclaw-gateway.service

[Service]
Type=simple
WorkingDirectory=/root/openclaw
ExecStart=/usr/bin/python3 /root/openclaw/agent_watchdog.py --daemon
Restart=always
RestartSec=10
Environment=PYTHONPATH=/root/openclaw

[Install]
WantedBy=multi-user.target
