# scripts/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway — AI Agency API Server
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/openclaw
ExecStart=/usr/bin/python3 gateway.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

# Hardening
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
